# ringboard/.config/systemd/user/ringboard-sync.service
[Unit]
Description=Sync Ringboard clipboard data to Cloudflare R2

[Service]
Type=oneshot
ExecStart=%h/.local/bin/ringboard-sync sync
