# ringboard/.config/systemd/user/ringboard-sync-pending.timer
[Unit]
Description=Check for pending Ringboard syncs

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
