# ringboard/.config/systemd/user/ringboard-sync-pending.service
[Unit]
Description=Process pending Ringboard syncs

[Service]
Type=oneshot
ExecStart=%h/.local/bin/ringboard-sync pending
